# budget-tracker.service
[Unit]
Description=Budget Tracker — personal budget tracker web app
After=network.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/scripts/start-standalone.sh
Restart=always
RestartSec=5
EnvironmentFile=__APP_DIR__/.env
Environment=PORT=__PORT__
Environment=HOSTNAME=0.0.0.0

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
